# deploy/go-librespot/go-librespot.service
[Unit]
Description=go-librespot (Spotify Connect) für Beocreate
After=network-online.target sound.target sigmatcpserver.service
Wants=network-online.target

[Service]
Type=simple
# Nicht als root: pi ist in der Gruppe audio und hat damit Zugriff auf die Soundkarte.
User=pi
ExecStart=/opt/go-librespot/go-librespot --config_dir /home/pi/.config/go-librespot
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
