# deploy/systemd/beocreate2.service
[Unit]
Description=Beocreate 2
After=network-online.target sigmatcpserver.service audiocontrol-shim.service
Wants=network-online.target audiocontrol-shim.service
Requires=sigmatcpserver.service
[Service]
Type=simple
User=root
WorkingDirectory=/opt/beocreate/beo-system
# Produktion: kein Debug-Flag. Für Diagnose "v" (Stufe 1) anhängen; NICHT "vv"
# (Stufe 2) — dessen Volume-Spam löst journald-Ratelimiting aus.
ExecStart=/usr/bin/node --require /opt/beocreate/beo-guard.js /opt/beocreate/beo-system/beo-server.js
Restart=on-failure
RestartSec=3
[Install]
WantedBy=multi-user.target
