# deploy/systemd/sigmatcpserver.service.d/override.conf
[Service]
ExecStart=
ExecStart=sigmatcpserver --alsa --enable-rest
