# deploy/audiocontrol-shim/audiocontrol-shim.service
[Unit]
Description=AudioControl2 shim for Beocreate 2 (sources/now-playing/volume API on :81)
After=network.target sigmatcpserver.service
Wants=sigmatcpserver.service

[Service]
Type=simple
# amixer (DSPVolume) needs access to the sound card - root like beocreate2/sigmatcpserver.
User=root
ExecStart=/usr/bin/node /opt/beocreate/audiocontrol-shim/shim.js
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
